# com.etf-tool.wind-download.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.etf-tool.wind-download</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/apangduo/.workbuddy/binaries/python/versions/3.13.12/bin/python3</string>
        <string>/Users/apangduo/WorkBuddy/Claw/etf-tool-mvp/download_wind_full_data.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/apangduo/WorkBuddy/Claw/etf-tool-mvp</string>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/apangduo/WorkBuddy/Claw/etf-tool-mvp/logs/wind_download_out.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/apangduo/WorkBuddy/Claw/etf-tool-mvp/logs/wind_download_err.log</string>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
